# packaging/systemd/piperd.user.service
[Unit]
Description=Piper agent (rootless dev instance)
After=docker.service network-online.target

[Service]
Type=simple
ExecStart=%h/.local/bin/piperd
Environment=PIPER_HTTP_ADDR=:8080
Environment=PIPER_HTTPS_ADDR=:8443
# Relocate the embedded Caddy admin off the default :2019 so a rootless dev
# instance can coexist with (or fail legibly against) a system piperd, which
# keeps the default admin port (#211).
Environment=PIPER_CADDY_ADMIN=http://127.0.0.1:2020
Environment=XDG_DATA_HOME=%h/.piper/piperd
Environment=XDG_CONFIG_HOME=%h/.piper/piperd
EnvironmentFile=-%h/.piper/piperd.env
Restart=on-failure
RestartSec=2s

[Install]
WantedBy=default.target
